# config/home/dot.config/systemd/user/check-battery.timer
[Unit]
Description=Timer to check battery status

[Timer]
OnCalendar=*:0/5
Unit=check-battery.service

[Install]
WantedBy=timers.target
