# config/home/dot.config/systemd/user/check-battery.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/sfusco/.config/systemd/user/check-battery.sh
